# sensors_network/project/LivoloLightsNode/src/ble/central_node/configs/systemd/ble_stats_panel.service
[Unit]
Description=Livolo BLE Stats Panel
After=network.target

[Service]
User=alarm
#Restart=on-failure
WorkingDirectory=/home/alarm/livolo/stats_panel/
ExecStart=/usr/bin/sudo /home/alarm/livolo/stats_panel/ble_stats_panel.py /home/alarm/livolo/configs/ble_stats_panel.conf

[Install]
WantedBy=multi-user.target
